# agent/nginx.service.d/nginx-manager-agent-recovery.conf
[Unit]
Requires=nginx-manager-agent-recover.service
After=nginx-manager-agent-recover.service
